# undervest.service.in
[Unit]
Description=Early fuzzy detection of unknown e-mail recipients daemon
After=syslog.target network.target

[Service]
Type=simple
Environment="UV_SOCKET=local:@RUNDIR@/undervest.sock"
Environment="UV_CONF=@SYSCONFDIR@/undervest.conf"
EnvironmentFile=-@SYSCONFDIR@/sysconfig/undervest
ExecStartPre=-/sbin/portrelease undervest
ExecStart=@SBINDIR@/undervest -n -c "$UV_CONF" $UV_EXTRA_FLAGS "$UV_SOCKET"
ExecReload=@BINDIR@/kill -USR1 ${MAINPID}
ExecStop=@BINDIR@/kill -TERM ${MAINPID}
StandardOutput=syslog
StandardError=syslog

[Install]
WantedBy=multi-user.target
